# roles/jupyterhub/templates/profile.d/local_bin.sh
{{ ansible_managed | comment() }}

export PATH=/usr/local/bin:$PATH